# files/light.service
[Unit]
Description=Action Light

[Service]
ExecStart=/bin/light.start

[Install]
WantedBy=multi-user.target
